# src/pre_processamento/crontab
*/10 * * * * cd /app && /usr/local/bin/python3 pre_processamento.py >> /var/log/script.log 2>&1
